# smclogweb.service
[Unit]
Description=SMC Binance Futures Trading Service
After=network-online.target mihomo.service smclogweb.service
Wants=network-online.target
Requires=mihomo.service
Wants=smclogweb.service

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/strategy
EnvironmentFile=-/home/ubuntu/strategy/smctrade.env
ExecStart=/home/ubuntu/strategy/venv/bin/python /home/ubuntu/strategy/main.py --config /home/ubuntu/strategy/config_enhanced.json
Restart=always
RestartSec=5
StandardOutput=append:/home/ubuntu/strategy/logs/smctrade-service.log
StandardError=append:/home/ubuntu/strategy/logs/smctrade-service.log

[Install]
WantedBy=multi-user.target
